# dotfiles/config/sv/astal-bar/run
#!/bin/bash

exec 2>&1
cd ~/.local/lib/astal-bar || exit 1
exec lua init.lua
