# dotfiles/config/sv/light-locker/run
#!/bin/sh

exec light-locker --late-locking
