# dotfiles/config/sv/slacksim/run
#!/bin/sh

sleep 2m
xdotool key "Super_L+1"
sleep 10s
xdotool key "Super_L+5"