# dotfiles/etc/sv/greetd/run
#!/bin/sh
# Needed because of pam_ssh needing large mlock
ulimit -l 65536
exec 2>&1
exec greetd
